# docker/dynamic-app/chatvia-app/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"
autorestart=true

[program:pm2]
command=pm2-runtime app.js
autorestart=true
